# etc/systemd/system/do_de.timer
[Unit]
Description=Run do_de hourly

[Timer]
OnCalendar=hourly
Unit=do_de.service

[Install]
WantedBy=timers.target
